# scripts/gateway-task-backup.xml
<?xml version="1.0" encoding="UTF-16"?>

<Task version="1.2" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">

  <RegistrationInfo>

    <Date>2026-03-17T20:26:18</Date>

    <Author>A5_PRO\IkeFl</Author>

    <URI>\OpenClaw Gateway</URI>

  </RegistrationInfo>

  <Principals>

    <Principal id="Author">

      <UserId>S-1-5-21-4149131261-2859786337-62391995-1001</UserId>

      <LogonType>InteractiveToken</LogonType>

    </Principal>

  </Principals>

  <Settings>

    <DisallowStartIfOnBatteries>true</DisallowStartIfOnBatteries>

    <StopIfGoingOnBatteries>true</StopIfGoingOnBatteries>

    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>

    <IdleSettings>

      <Duration>PT10M</Duration>

      <WaitTimeout>PT1H</WaitTimeout>

      <StopOnIdleEnd>true</StopOnIdleEnd>

      <RestartOnIdle>false</RestartOnIdle>

    </IdleSettings>

  </Settings>

  <Triggers>

    <LogonTrigger>

      <StartBoundary>2026-03-17T20:26:00</StartBoundary>

    </LogonTrigger>

  </Triggers>

  <Actions Context="Author">

    <Exec>

      <Command>C:\Users\IkeFl\.openclaw\gateway.cmd</Command>

    </Exec>

  </Actions>

</Task>
